# scripts/cron.hourly/regulate
#!/bin/sh -e
# clean downloaded and transcoded content
regulator
